# scripts/stoma-scheduler.service
[Unit]
Description=Stoma Research Intelligence Scheduler
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=cordlesssteve
Group=cordlesssteve
WorkingDirectory=/home/cordlesssteve/projects/Utility/RESEARCH/Stoma
Environment=PYTHONPATH=/home/cordlesssteve/projects/Utility/RESEARCH/Stoma
ExecStart=/usr/bin/python3 -m stoma.scheduler.manager
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stoma-scheduler

# Resource limits
LimitNOFILE=65536
MemoryLimit=1G

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/cordlesssteve/projects/Utility/RESEARCH/Stoma

[Install]
WantedBy=multi-user.target